# .config/systemd/user/backup_timer.service
[Unit]
Description=Restic Backup Service

[Service]
Type=simple
ExecStart=/home/arneb/.local/bin/backup

